# dotfiles/config/systemd/user/default.target.wants/gitwatch.service
[Unit]
Description=Gitwatch

[Service]
ExecStart=/usr/bin/gitwatch -r origin /home/tjark/Nextcloud/Notizen/

[Install]
WantedBy=default.target
